# docker/http-servers/http-supervisord/source-files/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/var/log/supervisor/supervisord.log

[program:nginx]
command=service nginx start

[program:fcgiwrap]
command=service fcgiwrap start 
